# contrib/launchd/com.mcprt.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcprt.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mcprt</string>
        <string>serve</string>
        <string>--config</string>
        <string>/Users/YOUR_USERNAME/.config/mcprt/mcprt.toml</string>
        <string>--log-level</string>
        <string>info</string>
    </array>
    <!-- mcprt starts all MCP servers on-demand; RunAtLoad is true for the
         daemon itself but the MCP child processes start only when a client
         connects. Total idle RSS: ~30 MB for the daemon, 0 MB for children. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/mcprt.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/mcprt.err</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
